# telia-zone.list
/usr/local/etc/zonewatcher.conf
/usr/local/bin/setbridge.sh
/etc/init.d/premiumd
/etc/rc.d/S98zone_watch
/usr/lib/libstdc++.so.6
/usr/lib/libstdc++.so.6.0.16-gdb.py
/usr/local/bin/hotspotzone.sh
/etc/nginx/wildcard.premiumzone.com.key
/etc/rc.d/S99zoneredird
/etc/nginx/wildcard.premiumzone.com.chain.pem
/etc/rc.d/S99premiumd
/usr/local/sbin/zonewatcher
/usr/local/etc/premiumd.conf
/usr/local/etc/zoneredirect.conf
/usr/local/sbin/zoneredir_daemon
/www/lua/zone_common.lua
/usr/local/sbin/zoneredir
/etc/init.d/zoneredird
/usr/lib/libstdc++.so
/usr/local/sbin/premium
/usr/local/sbin/zone_daemon
/usr/lib/libstdc++.so.6.0.16
/usr/local/bin/arpzone.sh
/usr/local/bin/serial.sh
/usr/local/bin/zone_watch
/www/lua/zone_deregister.lua
/etc/init.d/zone_watch
/usr/local/sbin/premiumd_daemon
/www/lua/zone_register.lua
